# templates/etc/init/cinder-start-loop.conf
description  "Set up loop device for cinder-volumes volume group"

start on local-filesystems
task
exec /usr/local/sbin/cinder-start-loop
